# ansible/roles/kiosk/templates/kiosk.service.j2
# Managed by Ansible - do not edit manually.
[Unit]
Description=Kiosk display - Grafana dashboard via cage + Chromium
After=multi-user.target snapd.seeded.service
Wants=multi-user.target

[Service]
User={{ kiosk_user }}
Group={{ kiosk_user }}
PAMName=login
TTYPath=/dev/tty7
Environment=XDG_RUNTIME_DIR=/run/user/{{ kiosk_uid }}
Environment=XDG_SESSION_TYPE=wayland
Environment=LIBSEAT_BACKEND=logind
ExecStartPre=+/usr/bin/chvt 7
ExecStartPre=+/bin/mkdir -p /run/user/{{ kiosk_uid }}
ExecStartPre=+/bin/chown {{ kiosk_user }}:{{ kiosk_user }} /run/user/{{ kiosk_uid }}
ExecStartPre=+/bin/chmod 700 /run/user/{{ kiosk_uid }}
ExecStart=/usr/bin/cage -- \
    /snap/bin/chromium \
    --no-sandbox \
    --kiosk \
    --noerrdialogs \
    --disable-translate \
    --no-first-run \
    --disable-infobars \
    --disable-session-crashed-bubble \
    {{ grafana_kiosk_url }}
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
